# root/etc/s6-overlay/s6-rc.d/slp2mp4-watcher-deps/run
#!/usr/bin/with-contenv bash
set -euo pipefail
mkdir -p /config/SlippiReplays /.config /config
mkdir -p /config/SlippiReplays/.work /config/SlippiReplays/.failed
